# files/etc_init.d_weblogic
#!/bin/bash
#
# sites    This shell script takes care of starting and stopping
#
# chkconfig: - 50 50
# description: sites

export MWA_ENV=DEV
export USER_MEM_ARGS="-Xmx2048m -XX:MaxPermSize=256m"
export PRE_CLASSPATH=/opt/oracle/sites/home/bin:/opt/oracle/weblogic/user_projects/applications/sites_domain/cs/WEB-INF/lib/commons-lang-2.4.jar
export PATH=/opt/oracle/sites/bin:$PATH
export JAVA_OPTIONS="-Dfile.encoding=UTF-8 -Dorg.owasp.esapi.resources=/opt/oracle/weblogic/user_projects/applications/sites_domain/cs/WEB-INF/classes -Duser.timezone=UTC -Dnet.sf.ehcache.enableShutdownHook=true -Djava.net.preferIPv4Stack=true -Djava.io.tmpdir=/var/tmp/weblogic"

case $1 in
start)
 # cd /opt/oracle/weblogic/user_projects/applications/sites_domain
cd /opt/oracle/weblogic/user_projects/domains/SitesDomain/
 su weblogic  -c './startWebLogic.sh 2>&1 >/var/log/weblogic/sites.log' &
;;
stop)
 # cd /opt/oracle/weblogic/user_projects/applications/sites_domain/bin
cd /opt/oracle/weblogic/user_projects/domains/SitesDomain/bin
 su weblogic -c './stopWebLogic.sh'
;;
run)
 # cd /opt/oracle/weblogic/user_projects/applications/sites_domain
cd /opt/oracle/weblogic/user_projects/domains/SitesDomain/
 su weblogic  -c './startWebLogic.sh'
;;
status)
        # PROCESS_STRING="^.*/opt/oracle/.*weblogic.NodeManager.*"
        PROCESS_STRING="\/bin\/java \-client.* weblogic\.Server$"
        SERVICE_NAME=`/bin/basename $0`

        OLDPID=`/usr/bin/pgrep -f "$PROCESS_STRING"`
        if [ "$OLDPID" != "" ]; then
            /bin/echo "$SERVICE_NAME is running (pid: $OLDPID)"
        else
            /bin/echo "$SERVICE_NAME is stopped"
        fi
        RETVAL=$?
        ;;

*) echo "usage: {start|stop|run|status}"
;;
esac



